# sirius-authenticate.timer
[Unit]
Description=Run Sirius Audio Authenticate daily at 3 AM

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target